# scripts/vboxtool/vboxtoolinit
#!/bin/sh

. /etc/rc.subr

name="vboxtool"
rcvar=`vboxtool_enable`
command="/usr/local/bin/${name}"
command_args="autostart"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"